# systemd/ngrok.service
[Unit]
Description=ngrok tunnel for BISINDO
After=network-online.target

[Service]
Type=simple
ExecStart=/home/kevin/.local/bin/ngrok http --region=us --no-tui --log=stdout 5000
Restart=on-failure
RestartSec=5
Environment=HOME=/home/kevin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
